# packaging/root/usr/share/init/dropship.conf
start on (filesystem and net-device-up)

chdir /var/lib/dropship
console log

script
	/usr/local/bin/dropship --config=/etc/dropship/dropship.json
end script
