# packaging/root/usr/share/dropship/init/dropship.conf
start on (filesystem and net-device-up)

console log

script
	sudo -u root /usr/local/bin/dropship agent --config=/etc/dropship.d/dropship.hcl
end script
